# 金元Ɉ/thoughts_bot/juno.service
[Unit]
Description=Juno Montana Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/ACP_1/金元Ɉ/thoughts_bot
ExecStart=/usr/bin/python3 -u /root/ACP_1/金元Ɉ/thoughts_bot/unified_bot.py
Restart=always
RestartSec=30
EnvironmentFile=/root/ACP_1/金元Ɉ/thoughts_bot/.env

[Install]
WantedBy=multi-user.target
